# deployment/systemd/podcastfy-frontend.service
[Unit]
Description=Podcastfy Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/podcastfy/frontend

# Environment
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/var/www/podcastfy/frontend/.env.local

# Start command
ExecStart=/usr/bin/npm start

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=podcastfy-frontend

[Install]
WantedBy=multi-user.target
